# units/autosocks-@.socket
[Unit]
Description=autosocks socket (:%I -> %J)

[Socket]
ListenStream=127.0.0.1:%i

[Install]
WantedBy=sockets.target
